# sesame-tts/supervisord.conf
[unix_http_server]
file=/var/run/supervisor.sock

[supervisord]
logfile=/var/log/supervisord.log
logfile_maxbytes=20MB
loglevel=info
pidfile=/var/run/supervisord.pid
nodaemon=false

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

# ── CSM GPU 0 (always on) ─────────────────────────────────────────────────────
[program:csm-0]
command=python /app/csm_server.py
directory=/app
autostart=true
autorestart=true
startretries=5
stdout_logfile=/var/log/csm-0.log
stderr_logfile=/var/log/csm-0.log
stdout_logfile_maxbytes=50MB
redirect_stderr=true
environment=
    HF_HOME="/model-cache",
    VOICE_REFS_DIR="/model-cache/voice-refs",
    PYTORCH_DISABLE_CUDA_GRAPHS="1",
    OMP_NUM_THREADS="1",
    MKL_NUM_THREADS="1",
    TOKENIZERS_PARALLELISM="false",
    CUDA_LAUNCH_BLOCKING="0",
    CUDA_VISIBLE_DEVICES="0",
    CSM_PORT="8080"

# ── CSM GPU 1 (dual-GPU mode only — enabled by setup_vast.sh DUAL_GPU=1) ──────
# To enable: supervisorctl update && supervisorctl start csm-1
[program:csm-1]
command=python /app/csm_server.py
directory=/app
autostart=false
autorestart=true
startretries=5
stdout_logfile=/var/log/csm-1.log
stderr_logfile=/var/log/csm-1.log
stdout_logfile_maxbytes=50MB
redirect_stderr=true
environment=
    HF_HOME="/model-cache",
    VOICE_REFS_DIR="/model-cache/voice-refs",
    PYTORCH_DISABLE_CUDA_GRAPHS="1",
    OMP_NUM_THREADS="1",
    MKL_NUM_THREADS="1",
    TOKENIZERS_PARALLELISM="false",
    CUDA_LAUNCH_BLOCKING="0",
    CUDA_VISIBLE_DEVICES="1",
    CSM_PORT="8081"

# ── nginx (dual-GPU mode only — enabled by setup_vast.sh DUAL_GPU=1) ──────────
[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=false
autorestart=true
startretries=3
stdout_logfile=/var/log/nginx.log
stderr_logfile=/var/log/nginx.log
redirect_stderr=true

# ── Cloudflare Tunnel ────────────────────────────────────────────────────────
[program:cloudflared]
command=/usr/local/bin/cloudflared tunnel --no-autoupdate run
autostart=false
autorestart=true
startretries=3
stdout_logfile=/var/log/cloudflared.log
stderr_logfile=/var/log/cloudflared.log
redirect_stderr=true
